# rootfs/etc/s6-overlay/s6-rc.d/init-killswitch/up
/etc/s6-overlay/scripts/init-killswitch.sh