# rootfs/etc/s6-overlay/s6-rc.d/init-environment/up
/etc/s6-overlay/scripts/init-environment.sh