# rootfs/etc/s6-overlay/s6-rc.d/init-firewall-pre/up
/etc/s6-overlay/scripts/init-firewall-pre.sh