# rootfs/etc/s6-overlay/s6-rc.d/svc-qbittorrent/finish
#!/bin/bash
echo "[tunnelvision] qBittorrent-nox exited with code $1, signal $2"
