# rootfs/etc/s6-overlay/s6-rc.d/svc-qbittorrent/run
#!/command/with-contenv bash
# ============================================================================
# TunnelVision — qBittorrent-nox Service
# Skips in setup mode — VPN must be configured first.
# ============================================================================

QBT_ENABLED=${QBT_ENABLED:-true}
if [ "$QBT_ENABLED" != "true" ]; then
    echo "[tunnelvision] qBittorrent disabled (QBT_ENABLED=false)"
    exec sleep infinity
fi

SETUP_REQUIRED=$(cat /var/run/tunnelvision/setup_required 2>/dev/null || echo "false")
if [ "$SETUP_REQUIRED" = "true" ]; then
    echo "[tunnelvision] Setup mode — qBittorrent waiting for VPN configuration"
    exec sleep infinity
fi

echo "[tunnelvision] Starting qBittorrent-nox..."

exec s6-setuidgid tunnelvision \
    qbittorrent-nox \
        --confirm-legal-notice \
        --profile=/config \
        --webui-port="${WEBUI_PORT:-8080}" \
        2>&1
