# rootfs/etc/s6-overlay/s6-rc.d/init-wireguard/up
/etc/s6-overlay/scripts/init-wireguard.sh